# home/dot_config/systemd/user/clamp-bt-volume.service
[Unit]
Description=Clamp Bluetooth volume to a maximum

[Service]
Type=oneshot
ExecStart=%h/.local/bin/clamp-bt-volume
