# home/dot_config/systemd/user/private_lid-suspend.service
[Unit]
Description=Suspend on lid close
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl suspend

[Install]
WantedBy=suspend.target
